# modules/logstash/templates/init_logstash-agent.conf.erb
# logstash agent
description "Logstash agent"

start on stopped rc RUNLEVEL=[345]
stop on runlevel [!2345]

respawn 
respawn limit 2 10

# run with lowest priority
nice 19
chdir /tmp

script
    JAR='/opt/logstash/logstash.jar'
    CONF='/etc/logstash/logstash-agent.conf'
    LOG='/var/log/logstash/agent.log'
    PLUGINPATH='/opt'

    exec sudo java -jar $JAR agent --pluginpath $PLUGINPATH -f $CONF -v > $LOG 2>&1
end script
