# modules/graylog2/templates/init_graylog2-server.conf.erb
# graylog2 server
description "Graylog2 server"

start on stopped rc RUNLEVEL=[345]
stop on runlevel [!2345]

respawn 
respawn limit 2 10

# run with lowest priority
nice 19

script
    JAR='<%= scope.lookupvar('graylog2::graylog_folder') -%>/graylog2-server/graylog2-server.jar'
    CONF='<%= scope.lookupvar('graylog2::graylog_folder') -%>/graylog2-server/graylog2.conf'
    LOG='<%= scope.lookupvar('graylog2::graylog_log_folder') -%>/server.log'
    PID='<%= scope.lookupvar('graylog2::graylog_folder') -%>/graylog2-server/graylog2.pid'

    exec sudo -i -u graylog java -jar $JAR -f $CONF -p $PID 2>&1 > $LOG
end script